# deploy/tokoonline.service
[Unit]
Description=Toko MDT (tokoonline) - Go HTTP server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=tokoonline
Group=tokoonline
WorkingDirectory=/opt/tokoonline
EnvironmentFile=/opt/tokoonline/.env
ExecStart=/opt/tokoonline/bin/server
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/tokoonline/uploads /opt/tokoonline/logs

[Install]
WantedBy=multi-user.target
